# BackgroundTaskHost.yml
Name: BackgroundTaskHost.exe
Description: Microsoft-signed Windows component used to host and execute background tasks associated with applications and system components. While not malicious, it can act as a trusted execution context for attacker-controlled background logic.
Aliases:
  - Alias: BackgroundTaskHost64.exe
Author: Prakash Thakur
Created: 2025-12-20
Commands:
  - Command: BackgroundTaskHost.exe
    Description: Executes registered background task logic within a trusted Windows host process.
    Usecase: Can be abused as a signed binary proxy execution context where attacker-controlled background tasks run under a benign Windows process, complicating detection and triage.
    Category: Execute
    Privileges: User
    MitreID: T1218
    OperatingSystem: Windows 10 1803, Windows 10 1909, Windows 11
    Tags:
      - Technique: Signed Binary Proxy Execution
      - Context: Execution Masquerading

  - Command: schtasks /create /tn UserSync /tr notepad.exe /sc onlogon /rl LIMITED
    Description: Creates a user-level scheduled task that later executes within the BackgroundTaskHost.exe context.
    Usecase: Demonstrates how non-administrative attackers can indirectly cause execution under BackgroundTaskHost.exe at user logon.
    Category: Persistence
    Privileges: User
    MitreID: T1053
    OperatingSystem: Windows 10 All
Full_Path:
  - Path: C:\Windows\System32\BackgroundTaskHost.exe
  - Path: C:\Windows\SysWOW64\BackgroundTaskHost.exe
Code_Sample:
  - Code: https://learn.microsoft.com/windows/uwp/launch-resume/support-your-app-with-background-tasks
Detection:
  - IOC: BackgroundTaskHost.exe spawning unexpected child processes
  - IOC: BackgroundTaskHost.exe executing shortly after user-level task creation
  - Analysis: https://strontic.github.io/xcyclopedia/library/backgroundTaskHost.exe-8B50BFD5811304543479B20D0A281C56.html
Resources:
  - Link: https://learn.microsoft.com/windows/uwp/launch-resume/background-tasks
  - Link: https://attack.mitre.org/techniques/T1218/
  - Link: https://strontic.github.io/xcyclopedia/library/backgroundTaskHost.exe-8B50BFD5811304543479B20D0A281C56.html
Acknowledgement:
  - Person: Prakash Munimsingh Thakur
    Handle: 'https://github.com/4renSick''

